# crontab

* * * * * (cd ~/GTD/; git add reference/linux_config; git commit -m "Autoupdate linux config submodule")
